# config/deploy/supervisor/conf.d/celery_beat.conf
# ================================
#  celery beat
# ================================

# the name of your supervisord program
[program:celery_beat]

# Set full path to celery program if using virtualenv
command=/opt/venv/billing_api/bin/newrelic-admin run-program /opt/venv/billing_api/bin/celery -A server.taskapp.celery:app beat -s /opt/apps/billing_api/celerybeat-schedule -l info

# The directory to your Django project
directory=/opt/apps/billing_api

# Process
priority=905
numprocs=1
process_name=%(program_name)s_%(process_num)d
startsecs=10
autostart=true
autorestart=true
stopwaitsecs=60
killasgroup=true
stopasgroup=true
stopsignal=KILL

# If supervisord is run as the root user, switch users to this UNIX user account
# before doing any processing.
user=api_user

# logs
redirect_stderr=true
stderr_logfile=/var/log/billing_api/celery_beat.err.log
stdout_logfile=/var/log/billing_api/celery_beat.out.log

